# etc/systemd/system/wall2notify.service
[Unit]
Description=NODE wall-to-notify-send service

[Service]
Type=simple
ExecStart=bash -ic "wall2notify"

[Install]
WantedBy=multi-user.target
